# config/systemd/user/dev-backup.timer
[Unit]
Description=Run dev backup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=45m

[Install]
WantedBy=timers.target
